# deploy/albatross-update-reboot.path
[Unit]
Description=Watch for Albatross update reboot requests
After=local-fs.target

[Path]
PathExists=/home/albatross/albatross/updates/reboot_requested
Unit=albatross-update-reboot.service

[Install]
WantedBy=multi-user.target
